# linux/config/containers/systemd/adguardhome.container
[Container]
Image=docker.io/adguard/adguardhome:latest
AutoUpdate=registry
Pod=adguardhome.pod
Volume=../adguardhome/conf:/opt/adguardhome/conf:Z
