# linux/config/containers/systemd/flaresolverr.container
[Container]
Image=ghcr.io/flaresolverr/flaresolverr:latest
AutoUpdate=registry
Pod=flaresolverr.pod
EnvironmentFile=../flaresolverr/flaresolverr.env
