# linux/config/containers/systemd/librespeed.container
[Container]
Image=lscr.io/linuxserver/librespeed:latest
AutoUpdate=registry
Pod=librespeed.pod
