# linux/config/containers/systemd/openspeedtest.container
[Container]
Image=docker.io/openspeedtest/latest:latest
AutoUpdate=registry
Pod=openspeedtest.pod
