# linux/config/containers/systemd/suwayomi.container
[Container]
Image=ghcr.io/suwayomi/tachidesk:latest
AutoUpdate=registry
Pod=suwayomi.pod
EnvironmentFile=../suwayomi/suwayomi.env
