# linux/config/containers/systemd/syncyomi.container
[Container]
Image=ghcr.io/syncyomi/syncyomi:latest
AutoUpdate=registry
Pod=syncyomi.pod
Volume=../syncyomi/config/config.toml:/config/config.toml:Z
